# adwi/infra/launchd/com.suneel.adwi-autoresearch-night.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.suneel.adwi-autoresearch-night</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/MAC/SuneelWorkSpace/adwi/bin/adwi-autoresearch-night</string>
        <string>2</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/MAC/SuneelWorkSpace/logs/autoresearch/adwi/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/MAC/SuneelWorkSpace/logs/autoresearch/adwi/launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/MAC</string>
        <key>TERM</key>
        <string>xterm-256color</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
